# services/level-display-py.service
[Unit]
Description=Level Display (Python + Framebuffer)
After=multi-user.target
Conflicts=level-sdl2.service sdl2-dashboard.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/embedded-linux/apps/level_display.py
Restart=on-failure
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
